# assets/wayclip-daemon.service
[Unit]
Description=Wayclip Daemon
After=graphical.target
Wants=graphical.target

[Service]
ExecStart=__WAYCLIP_DAEMON_PATH__
Restart=always
RestartSec=3

Environment=XDG_RUNTIME_DIR=%t
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=%t/bus
PassEnvironment=DISPLAY WAYLAND_DISPLAY XAUTHORITY

[Install]
WantedBy=default.target
